# systemd/kernel-proc-lab-collector.service
[Unit]
Description=Kernel Proc Lab JSONL collector
Wants=kernel-proc-lab.service
After=kernel-proc-lab.service

[Service]
Type=simple
ExecStart=/usr/bin/kernel-lab collector --output /var/log/kernel-proc-lab/events.jsonl
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
